# aios/services/updated/units/aios-updated.service
[Unit]
Description=AIOS Update and Recovery Service
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/libexec/aios/updated
Restart=on-failure

[Install]
WantedBy=multi-user.target
