# aios/services/policyd/units/aios-policyd.service
[Unit]
Description=AIOS Policy Engine

[Service]
Type=simple
ExecStart=/usr/libexec/aios/policyd
Restart=on-failure

[Install]
WantedBy=multi-user.target
