# aios/image/recovery/aios-recovery.target
[Unit]
Description=AIOS Recovery Target
Wants=aios-recovery-shell.service
After=basic.target
AllowIsolate=yes
