# aios/image/mkosi.extra/usr/lib/systemd/system/aios-runtimed.service
[Unit]
Description=AIOS Runtime Scheduler
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/aios/runtime/platform.env
ExecStart=/usr/libexec/aios/runtimed
Restart=on-failure

[Install]
WantedBy=multi-user.target
