# aios/image/mkosi.extra/usr/lib/systemd/system/aios-installer.target
[Unit]
Description=AIOS Installer Target
Wants=aios-installer.service
After=basic.target local-fs.target
AllowIsolate=yes
