# aios/image/mkosi.extra/usr/lib/systemd/system/aios-boot-evidence.service
[Unit]
Description=AIOS Hardware Boot Evidence Collector
After=multi-user.target aios-updated.service
Wants=aios-updated.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/libexec/aios/aios-boot-evidence.sh

[Install]
WantedBy=multi-user.target
