# aios/services/deviced/units/aios-deviced.service
[Unit]
Description=AIOS Device and Multimodal Service

[Service]
Type=simple
ExecStart=/usr/libexec/aios/deviced
Restart=on-failure

[Install]
WantedBy=multi-user.target
