# aios/image/firstboot/aios-ai-onboarding.service
[Unit]
Description=AIOS AI Onboarding Initialization
ConditionPathExists=!/var/lib/aios/onboarding/ai-onboarding-initialized
After=aios-firstboot.service local-fs.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/aios/firstboot/aios-firstboot.env
ExecStart=/usr/bin/bash /usr/libexec/aios/aios-ai-onboarding.sh
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
